# website_price_tracker/rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv bash
echo "[addon] stopping Website Price Tracker"
