# website_price_tracker/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
echo "[addon] starting Website Price Tracker (venv)"
exec /opt/venv/bin/python /app/app.py
