# jura_s8_price_addon/rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv bash
echo "[addon] stopping Jura S8 Price Tracker"
