# jura_s8_price_addon/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
echo "[addon] starting Jura S8 Price Tracker"
exec python3 /app/app.py
